# scripts/systemd/lh-hl-backfill-monitor.timer
[Unit]
Description=Hourly LiquidationHeatmap Hyperliquid Backfill Batch Monitor

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=120
Unit=lh-hl-backfill-monitor.service

[Install]
WantedBy=timers.target
