# scripts/systemd/lh-ingestion.timer
[Unit]
Description=Daily LiquidationHeatmap DuckDB Ingestion
After=binance-sync.timer

[Timer]
# Run at 03:30 daily (UTC) — 30 min after binance-sync starts at 03:00
OnCalendar=*-*-* 03:30:00
# If system was off, run on next boot
Persistent=true
# Randomize start within 5 minutes
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
